# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
nodaemon=true
logfile=/app/log/supervisor.log

[program:web]
command=bundle exec puma -e production -t 16:16 -p 80 --preload

[program:worker]
command=bundle exec sidekiq -C config/sidekiq.yml -e production

[program:cron]
command=cron -f
